# crontab
# START CRON JOB
* * * * * cd /app ; python3 src/ams_meervaart.py
# END CRON JOB

# Debug
# crond -f -l 2 -L /dev/stdout