# crontab
# START CRON JOB
MAILTO="alegomes@gmail.com"
*/10 * * * * cd /app ; python3 src/main.py
# END CRON JOB

# Debug
# crond -f -l 2 -L /dev/stdout